# files/supervisord-memcached.conf
[program:memcached]
command=memcached
autostart=true
autorestart=true
startsecs=2
user=memcache
